# 33.conf
[program:chain33]
command=/www/chain33/chain33
